# systemd/mirrord.service
[Unit]
Description=git repository mirror daemon

[Service]
ExecStart=/usr/bin/mirrord
Restart=always

[Install]
WantedBy=multi-user.target
